# overlay/system/usr/lib/systemd/system/hfd-service.service.d/wait-for-service.conf
[Service]
ExecStartPre=/bin/bash -c 'while [ "$(getprop init.svc.vendor.qti.vibrator)" != "running" ]; do sleep 0.5; done'
